# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=service nginx start

[program:php-fpm]
command= service php7.2-fpm start
